# supervisor.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
logfile=/var/log/teknofest/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=false
minfds=65536
minprocs=200
user=root
childlogdir=/var/log/teknofest

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ==========================================
# TEKNOFEST API - Main Application
# ==========================================
[program:teknofest-api]
command=/opt/teknofest-api/venv/bin/gunicorn src.app:app -c /opt/teknofest-api/gunicorn_config.py
directory=/opt/teknofest-api
user=teknofest
group=teknofest
numprocs=1
stdout_logfile=/var/log/teknofest/api-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/teknofest/api-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=30
killasgroup=true
stopasgroup=true
priority=999

environment=
    APP_ENV="production",
    PYTHONPATH="/opt/teknofest-api",
    API_WORKERS="%(ENV_API_WORKERS)s",
    AUTO_SCALE_ENABLED="true",
    WORKER_HEALTH_MONITORING="true",
    PROMETHEUS_ENABLED="true",
    DATABASE_URL="%(ENV_DATABASE_URL)s",
    REDIS_URL="%(ENV_REDIS_URL)s",
    SECRET_KEY="%(ENV_SECRET_KEY)s",
    HUGGING_FACE_HUB_TOKEN="%(ENV_HUGGING_FACE_HUB_TOKEN)s"

# ==========================================
# Worker Manager Service
# ==========================================
[program:teknofest-worker-manager]
command=/opt/teknofest-api/venv/bin/python /opt/teknofest-api/src/worker_manager.py --config /etc/teknofest/worker-config.json
directory=/opt/teknofest-api
user=teknofest
group=teknofest
numprocs=1
stdout_logfile=/var/log/teknofest/worker-manager-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/teknofest/worker-manager-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=30
priority=998

environment=
    APP_ENV="production",
    PYTHONPATH="/opt/teknofest-api",
    MIN_WORKERS="2",
    MAX_WORKERS="16",
    AUTO_SCALE_ENABLED="true",
    HEALTH_CHECKS_ENABLED="true"

# ==========================================
# Celery Workers (for async tasks)
# ==========================================
[program:teknofest-celery-worker]
command=/opt/teknofest-api/venv/bin/celery -A src.celery_app worker --loglevel=info --concurrency=4 -n worker-%%i@%%h
directory=/opt/teknofest-api
user=teknofest
group=teknofest
numprocs=2
numprocs_start=1
stdout_logfile=/var/log/teknofest/celery-worker-%%i-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/teknofest/celery-worker-%%i-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=30
killasgroup=true
stopasgroup=true
priority=997

environment=
    APP_ENV="production",
    PYTHONPATH="/opt/teknofest-api",
    CELERY_BROKER_URL="%(ENV_REDIS_URL)s",
    CELERY_RESULT_BACKEND="%(ENV_REDIS_URL)s"

# ==========================================
# Celery Beat (for scheduled tasks)
# ==========================================
[program:teknofest-celery-beat]
command=/opt/teknofest-api/venv/bin/celery -A src.celery_app beat --loglevel=info --pidfile=/var/run/celery-beat.pid
directory=/opt/teknofest-api
user=teknofest
group=teknofest
numprocs=1
stdout_logfile=/var/log/teknofest/celery-beat-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/teknofest/celery-beat-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=30
priority=996

environment=
    APP_ENV="production",
    PYTHONPATH="/opt/teknofest-api",
    CELERY_BROKER_URL="%(ENV_REDIS_URL)s",
    CELERY_RESULT_BACKEND="%(ENV_REDIS_URL)s"

# ==========================================
# Flower (Celery monitoring)
# ==========================================
[program:teknofest-flower]
command=/opt/teknofest-api/venv/bin/celery -A src.celery_app flower --port=5555 --basic_auth=admin:%(ENV_FLOWER_PASSWORD)s
directory=/opt/teknofest-api
user=teknofest
group=teknofest
numprocs=1
stdout_logfile=/var/log/teknofest/flower-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/teknofest/flower-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=10
priority=995

environment=
    APP_ENV="production",
    PYTHONPATH="/opt/teknofest-api",
    CELERY_BROKER_URL="%(ENV_REDIS_URL)s",
    CELERY_RESULT_BACKEND="%(ENV_REDIS_URL)s"

# ==========================================
# MCP Server
# ==========================================
[program:teknofest-mcp-server]
command=/opt/teknofest-api/venv/bin/python /opt/teknofest-api/src/mcp_server/production_server.py
directory=/opt/teknofest-api
user=teknofest
group=teknofest
numprocs=1
stdout_logfile=/var/log/teknofest/mcp-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/teknofest/mcp-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=30
priority=994

environment=
    APP_ENV="production",
    PYTHONPATH="/opt/teknofest-api",
    API_URL="http://localhost:8000",
    HUGGING_FACE_HUB_TOKEN="%(ENV_HUGGING_FACE_HUB_TOKEN)s",
    MODEL_NAME="%(ENV_MODEL_NAME)s"

# ==========================================
# Health Monitor
# ==========================================
[program:teknofest-health-monitor]
command=/opt/teknofest-api/venv/bin/python /opt/teknofest-api/scripts/health_monitor.py
directory=/opt/teknofest-api
user=teknofest
group=teknofest
numprocs=1
stdout_logfile=/var/log/teknofest/health-monitor-stdout.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=/var/log/teknofest/health-monitor-stderr.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
priority=993

environment=
    APP_ENV="production",
    PYTHONPATH="/opt/teknofest-api",
    HEALTH_CHECK_URL="http://localhost:8000/health",
    ALERT_WEBHOOK="%(ENV_ALERT_WEBHOOK)s"

# ==========================================
# Log Aggregator
# ==========================================
[program:teknofest-log-aggregator]
command=/opt/teknofest-api/venv/bin/python /opt/teknofest-api/scripts/log_aggregator.py
directory=/opt/teknofest-api
user=teknofest
group=teknofest
numprocs=1
stdout_logfile=/var/log/teknofest/log-aggregator-stdout.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=/var/log/teknofest/log-aggregator-stderr.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
priority=992

environment=
    APP_ENV="production",
    PYTHONPATH="/opt/teknofest-api",
    LOG_DIR="/var/log/teknofest",
    ELASTICSEARCH_URL="%(ENV_ELASTICSEARCH_URL)s"

# ==========================================
# Process Groups
# ==========================================
[group:teknofest-core]
programs=teknofest-api,teknofest-worker-manager,teknofest-mcp-server
priority=999

[group:teknofest-async]
programs=teknofest-celery-worker,teknofest-celery-beat,teknofest-flower
priority=998

[group:teknofest-monitoring]
programs=teknofest-health-monitor,teknofest-log-aggregator
priority=997

# ==========================================
# Event Listeners
# ==========================================
[eventlistener:teknofest-crashmail]
command=/opt/teknofest-api/venv/bin/crashmail -p teknofest -m admin@teknofest.com
events=PROCESS_STATE_EXITED
stdout_logfile=/var/log/teknofest/crashmail.log
stderr_logfile=/var/log/teknofest/crashmail-error.log

[eventlistener:teknofest-memmon]
command=/opt/teknofest-api/venv/bin/memmon -p teknofest-api=2GB -p teknofest-celery-worker=1GB -m admin@teknofest.com
events=TICK_60
stdout_logfile=/var/log/teknofest/memmon.log
stderr_logfile=/var/log/teknofest/memmon-error.log